# workers/deploy/cf-dispatcher.service
[Unit]
Description=cf dispatcher (chequeo barato -> encola jobs de corrección)
After=network.target docker.service
Wants=docker.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/clausina/core/workers
Environment=HOME=/root
Environment=REDIS_URL=redis://localhost:6380/0
ExecStart=/usr/bin/python3 /root/clausina/core/workers/dispatcher.py
